# deploy/zte-watchdog.service
[Unit]
Description=ZTE MC801A connection watchdog
After=network-online.target
Wants=network-online.target

# install.sh rewrites User= (to the installing user) and ExecStart= (to the
# repo's .venv python). Edit both by hand if you install this unit yourself.
[Service]
User=wojtek
EnvironmentFile=/etc/zte-watchdog.env
ExecStart=/usr/bin/python3 -m zte_watchdog --config /etc/zte-watchdog/config.toml
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
